# files/supervisor-conf/supervisor-worker-template.conf
[program:{{WORKER_NAME}}]
process_name=%(program_name)s_%(process_num)02d
command=php {{ARTISAN}} queue:work --timeout=0 --sleep=10
autostart=true
autorestart=true
user={{WEBSERVER_USER}}
redirect_stderr=true
stdout_logfile={{WORKER_LOG}}
